# macos/com.opencode.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.opencode.server</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/sh</string>
      <string>-c</string>
      <string>exec /Users/netanel/.bun/bin/bun /Users/netanel/.bun/bin/opencode serve --hostname 0.0.0.0 --port 4096</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>/Users/netanel</string>

    <key>StandardOutPath</key>
    <string>/Users/netanel/Library/Logs/opencode-server.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/netanel/Library/Logs/opencode-server.err.log</string>
  </dict>
</plist>
